# deploy/tecnipro-daily.service
[Unit]
Description=Tecnipro Daily Pipeline (Scrape + Process + Reports)
After=network.target

[Service]
Type=oneshot
User=tecnipro
Group=tecnipro
WorkingDirectory=/home/tecnipro/tecnipro-reportes
Environment="PATH=/home/tecnipro/tecnipro-reportes/venv/bin"
Environment="TZ=America/Santiago"
EnvironmentFile=/home/tecnipro/tecnipro-reportes/.env
ExecStart=/bin/bash /home/tecnipro/tecnipro-reportes/scripts/run_daily_production.sh
TimeoutStartSec=1200
StandardOutput=append:/var/log/tecnipro/daily.log
StandardError=append:/var/log/tecnipro/daily-error.log
